# ipkg/nodogsplash.conffiles
/etc/nodogsplash/nodogsplash.conf
/etc/nodogsplash/htdocs/splash.html
/etc/init.d/S65nodogsplash

